# s5.service
[Unit]
Description=S5 Server
 
[Service]
Type=oneshot
ExecStart=/etc/rc.d/init.d/ss5  -k stop
 
[Install]
WantedBy=multi-user.target

